# agentlabjournal-source-triage.timer
[Unit]
Description=Review AI source candidates twice daily

[Timer]
OnCalendar=*-*-* 06:45:00 Europe/Moscow
OnCalendar=*-*-* 18:45:00 Europe/Moscow
Persistent=true
Unit=agentlabjournal-source-triage.service

[Install]
WantedBy=timers.target
