# agentlabjournal-article-cycle.timer
[Unit]
Description=Generate Agent Lab Journal articles from the approved queue

[Timer]
OnCalendar=*-*-* 08,12,16,20:00:00 Europe/Moscow
Persistent=true
Unit=agentlabjournal-article-cycle.service

[Install]
WantedBy=timers.target
