# agentlabjournal-article-cycle.service
[Unit]
Description=Agent Lab Journal article generation cycle
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/agentlabjournal
EnvironmentFile=/root/.config/midas-telegram-bot.env
EnvironmentFile=/root/.config/agentlabjournal-dev.env
Environment=GIT_SSH_COMMAND=ssh -i /root/.ssh/grifun_github -o IdentitiesOnly=yes
ExecStart=/usr/bin/python3 /root/agentlabjournal/scripts/run-article-cycle.py
