# agentlabjournal-source-collector.service
[Unit]
Description=Collect AI article sources for Agent Lab Journal
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/agentlabjournal
ExecStart=/usr/bin/python3 /root/agentlabjournal/scripts/collect-article-sources.py
