# agentlabjournal-source-triage.service
[Unit]
Description=Review Agent Lab Journal source candidates
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/agentlabjournal
ExecStart=/usr/bin/python3 /root/agentlabjournal/scripts/triage-article-sources.py
