# agentlabjournal-source-collector.timer
[Unit]
Description=Collect AI article sources twice daily

[Timer]
OnCalendar=*-*-* 06:30:00 Europe/Moscow
OnCalendar=*-*-* 18:30:00 Europe/Moscow
Persistent=true
Unit=agentlabjournal-source-collector.service

[Install]
WantedBy=timers.target
